# files/etc/hotplug.d/iface/50-wireless-uplink
if [ "$INTERFACE" = 'wwan' ]; then
   /lib/wifi-uplink/wifi-darmstadt/captive.sh
fi
