# wifi-darmstadt-uplink/files/etc/hotplug.d/iface/99-wireless-uplink
if [ "$INTERFACE" = 'wwan' ]; then
   /lib/gluon/wan-dnsmasq/update.lua
   gluon-wan /lib/wifi-uplink/wifi-darmstadt/captive.sh
fi
